# overlay-lakitu/sys-apps/system-sysdaemons/files/system-sysdaemons.slice
[Unit]
Description=Slice for System Daemons
Before=slices.target
