# overlay-lakitu/chromeos-base/cloud-audit-config/files/cloud-audit-setup.service
[Unit]
Description=Setup audit rules
After=iptables-setup.service

[Service]
Type=oneshot
ExecStart=/usr/share/cloud/audit/setup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
