# project-moblab/chromeos-base/chromeos-bsp-moblab/files/init/mobmonitor-init.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Launch the Mob* Monitor for moblab"
author        "chromium-os-dev@chromium.org"

start on started boot-services
respawn

env LOG_DIR=/var/log/mobmonitor

pre-start script
  mkdir -p "${LOG_DIR}"
  chown -R moblab:moblab "${LOG_DIR}"
end script

exec python /usr/local/moblab/mobmonitor/mobmonitor.py
