# overlay-lakitu/chromeos-base/chromeos-init-systemd/files/dev-shm-remount.service
[Unit]
Description=Mount /dev/shm with 'noexec'
DefaultDependencies=false
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/mount -o remount,noexec /dev/shm
