# overlay-lakitu/chromeos-base/chromeos-init-systemd/files/mnt-stateful_partition-make-private.service
[Unit]
Description=Make /mnt/stateful_partition private
DefaultDependencies=false
After=mnt-stateful_partition.mount
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/mount --make-private /mnt/stateful_partition
