# project-mobbuild/chromeos-base/chromeos-bsp-mobbuild/files/init/mobbuild-init-begin.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on started boot-services
# This is the abstract job that marks the beginning of all mobbuild specific
# upstart jobs. All mobbuild jobs should start when this job is 'started'.

env UPSTART_LOG_DIR=/var/log/upstart
export UPSTART_LOG_DIR

pre-start script
  # TODO(pprabhu) Make this a feature of CrOS proper, and then remove this hack.
  # Only the user running this script can write to the folder. This is OK
  # because other upstart jobs should redirect their stdout/stderr from the job
  # file (i.e., the user writing logs is the same as the one running this job).
  if ! mkdir -m 0755 -p "${UPSTART_LOG_DIR}"; then
    logger -t "${UPSTART_JOB}" user.err \
        "Could not create log directory for my children. Bailing."
  fi
  initctl log-priority debug
end script
