# project-mobbuild/chromeos-base/chromeos-bsp-mobbuild/files/init/mobbuild-build-disk-init.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# event:mobbuild-build-disk-attached is generated by udev when the disk to be
# used as workspace is attached. See 65-mobbuild-build-disk-attached.rules
start on started mobbuild-init-begin and mobbuild-build-disk-attached

env MOBBUILD_BUILD_DIR=/b
export MOBBUILD_BUILD_DIR

script
  info () {
    echo "INFO: $@"
    logger -t "${UPSTART_JOB}" "$@"
  }
  err () {
    echo "ERR : $@"
    logger -t "${UPSTART_JOB}" -p user.err "$@"
  }

  # Redirect stdout, stderr to well known log files.
  exec >"${UPSTART_LOG_DIR}/${UPSTART_JOB}.log" 2>&1

  disk_drive="$(readlink -f /dev/disk/by-label/MOBBUILD_BUILD)"
  if [ -e "${disk_drive}" ]; then
    info "Mounting mobbuild build dir using ${disk_drive}"
  else
    err "No disk labeled MOBBUILD_BUILD found!"
    err "Mobbuild cannot be initialized."
    exit 1
  fi

  # Repeated executions of this script should succeed, causing a simple remount.
  umount "${MOBBUILD_BUILD_DIR}" || true
  if ! mount -o rw,exec,suid "${disk_drive}" "${MOBBUILD_BUILD_DIR}"; then
    err "Could not mount ${disk_drive} with required permissions"
    err "Mobbuild can not be initialized."
    exit 1
  fi

  # We don't change ownership recursively. We assume that files under this
  # folder are created with the right ownership.
  chown mobbuild:mobbuild "${MOBBUILD_BUILD_DIR}"
end script
