# project-mobbuild/chromeos-base/chromeos-bsp-mobbuild/files/init/mobbuild-dnsmasq.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Run dnsmasq to serve user supplied hosts"
author "chromium-os-dev@chromium.org"

# Don't block network-services so that creds disk loading failure doesn't block
# shill startup. shill startup failure is a much more serious offence.
start on started network-services and stopped mobbuild-creds-disk-init
stop on starting pre-shutdown
respawn

env MOBBUILD_HOME=/home/mobbuild

expect daemon
exec dnsmasq --no-dhcp-interface --local-service -u dnsmasq -g dnsmasq -i lo \
    --addn-hosts "${MOBBUILD_HOME}/hosts"
