# project-termina/chromeos-base/chromeos-bsp-termina/files/container-mount.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Container rootfs mounting service"
author          "chromium-os-dev@chromium.org"

start on started boot-services
task

script
  # The rootfs is expected to be the second block device, but fall back to 9p.
  if [ -b /dev/vdb ]; then
    mount /dev/vdb /mnt/container_rootfs
  else
    mount -t 9p -o trans=virtio,version=9p2000.L \
          container_rootfs /mnt/container_rootfs
  fi
end script
