# project-termina/chromeos-base/chromeos-bsp-termina/files/container.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Container launcher service"
author          "chromium-os-dev@chromium.org"

start on started container-mount
stop on stopping boot-services
task

exec run_oci -p chronos_containers /mnt/container_rootfs
