# overlay-oak/chromeos-base/chromeos-bsp-oak/files/99-mwifiex-wowlan.rules
KERNEL=="mlan*", SUBSYSTEM=="net", ACTION=="add", RUN+="/usr/sbin/ethtool -s %k wol u"
