# overlay-lakitu/chromeos-base/chromeos-init-systemd-dev/files/usr-local-remount.service
[Unit]
Description=Mount /usr/local with 'exec,suid'
DefaultDependencies=false
After=usr-local.mount
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/mount -o remount,exec,suid /usr/local
