# overlay-lakitu/sys-apps/systemd/files/systemd_upstart.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on starting boot-services

pre-start script
  mkdir -p /var/lib/dbus
  dbus-uuidgen --ensure

  mkdir --mode=2755 -p /var/log/journal
  chgrp systemd-journal /var/log/journal
end script

exec /usr/lib/systemd/systemd --system

post-stop script
  # The system is in inconsistent state if systemd crashes. If systemd was PID 1,
  # it would result in system restart. So simulate the same behavior here.
  logger -t "$UPSTART_JOB" "systemd exited unexpectedly; rebooting"
  uptime_sec=$(awk '{ print $1 }' </proc/uptime)
  metrics_client -t Uptime.SystemdCrash "${uptime_sec}" 1 1000000 50
  reboot
end script
